# deploy/coding-plan-mask.service
[Unit]
Description=Coding Plan Mask - 本地代理转发工具
Documentation=https://github.com/systemime/coding-plan-mask
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/project/coding-plan-mask
ExecStart=/opt/project/coding-plan-mask/bin/coding-plan-mask -config /opt/project/coding-plan-mask/config/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# 环境变量
Environment="GOMAXPROCS=2"

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

# 安全加固
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/project/coding-plan-mask
ReadOnlyPaths=/etc

# 标准输出和错误
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coding-plan-mask

# 优雅关闭
TimeoutStopSec=30s
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
